# launchd/com.polybot.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polybot.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.11</string>
        <string>/Users/brij/polybot/scripts/dashboard_server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/brij/polybot</string>

    <!-- Keep alive — this is a persistent server -->
    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/brij/polybot/logs/dashboard.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brij/polybot/logs/dashboard.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/brij</string>
        <key>DASHBOARD_PORT</key>
        <string>8767</string>
    </dict>
</dict>
</plist>
