# launchd/com.polybot-interactions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polybot-interactions</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/polybot &amp;&amp; source .venv/bin/activate &amp;&amp; .venv/bin/python core/discord_interactions.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/brij/polybot/.venv/bin</string>
        <key>DB_PATH</key>
        <string>/Users/brij/polybot/data/polybot.db</string>
        <key>INTERACTIONS_PORT</key>
        <string>8765</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/brij/polybot</string>

    <!-- Keep alive — restart immediately if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 5s before restarting on crash -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/brij/polybot/logs/interactions.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brij/polybot/logs/interactions.log</string>
</dict>
</plist>
